# Powernode Firmware/rootfs/etc/udev/rules.d/65-automount.rules
KERNEL=="sd[a-z]*", RUN+="/etc/udev/rules.d/automount.sh $ATTR{idVendor} $ATTR{idProduct}"
