# Powernode Firmware/rootfs/etc/udev/rules.d/70-cdrom.rules
ACTION=="add", KERNEL=="sr[0-9]*", SYMLINK+="cdrom", RUN+="/etc/udev/rules.d/automount.sh"